# debian/pi-wifi-setup.timer
[Timer]
OnActiveSec=15

[Install]
WantedBy=timers.target
